# releng/airootfs/etc/systemd/system/copyfiles.service
[Unit]
After=systemd-udev-settle.service

[Service]
ExecStart=/usr/local/bin/copyfiles.sh

[Install]
WantedBy=multi-user.target
