# releng/airootfs/etc/systemd/system/setfont.service
[Unit]
After=livecd-talk.service
ConditionKernelCommandLine=accessibility=on

[Service]
ExecStart=/usr/local/bin/setfont.sh

[Install]
WantedBy=multi-user.target
